# dist/orchestrator.service
[Unit]
Description=Orchestrator

[Service]
Type=simple
WorkingDirectory=%h/orchestrator
ExecStart=%h/orchestrator/orchestrator daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
TimeoutStopSec=3600
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
Environment=HOME=%h

[Install]
WantedBy=default.target
